# launchd/com.sportsbrain.auto-heal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sportsbrain.auto-heal</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/philiprassillier/sportsbrain/scripts/auto_heal_cron.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/philiprassillier/sportsbrain</string>

    <!-- Every 10 minutes -->
    <key>StartInterval</key>
    <integer>600</integer>

    <key>StandardOutPath</key>
    <string>/Users/philiprassillier/sportsbrain/results/launchd_auto_heal.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/philiprassillier/sportsbrain/results/launchd_auto_heal.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
