# launchd/com.sportsbrain.auto-retrain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sportsbrain.auto-retrain</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/philiprassillier/sportsbrain/scripts/auto_retrain_cron.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/philiprassillier/sportsbrain</string>

    <!-- Fires twice per day: 06:00 and 18:00 (every 12h during WM 2026) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/philiprassillier/sportsbrain/results/launchd_auto_retrain.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/philiprassillier/sportsbrain/results/launchd_auto_retrain.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
